# v2muser/runit/v2muser/run
#!/bin/sh
exec python3 /usr/bin/v2muser/v2muser.py
