# nginx-status/runit/srvstatus/run
#!/bin/sh
exec python3 /usr/bin/srvstatus/client-linux.py 2>&1
